# deploy/configs/sentinel.service
[Unit]
Description=Project Sentinel Military Alert Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sentinel
Group=sentinel
WorkingDirectory=/home/deploy/sentinel
ExecStart=/home/deploy/sentinel/venv/bin/python sentinel.py --config /etc/sentinel/config.yaml
Restart=always
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=300

# Environment -- loaded by PID 1 (root) before dropping to User=sentinel,
# so the sentinel user never reads the file directly
EnvironmentFile=/etc/sentinel/sentinel.env

# --- Security sandbox ---
# Prevent privilege escalation
NoNewPrivileges=yes

# Filesystem isolation
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/sentinel /var/log/sentinel
PrivateTmp=yes

# Device and kernel isolation
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Restrict capabilities and syscalls
CapabilityBoundingSet=
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel

[Install]
WantedBy=multi-user.target
